# main/database/bybit/crontab
*/1 *   * * *   ubuntu  bash ${HOMETRADE}/main/database/bybit/monitor.sh
0   1   * * *   ubuntu  nohup python ${HOMETRADE}/main/database/bybit/download_pastdata.py --fr `date "+\%Y\%m\%d" --date '2 day ago'` --to `date "+\%Y\%m\%d" --date '1 day ago'` --num 10000 --fn trade --fname bybit.gz --update > ${HOMETRADE}/main/log/bybit_download.`date "+\%Y\%m\%d"`.log 2>&1 &
