# main/database/bitflyer/crontab
*/1 *   * * *   ubuntu  bash /home/ubuntu/kktrade/main/database/bitflyer/monitor.sh
0   3   * * *   ubuntu  sleep 10 && pkill python