# main/database/bitflyer/crontab
*/1 *   * * *   ubuntu  bash /home/ubuntu/kktrade/main/database/bitflyer/monitor.sh
