# main/database/binance/crontab
*/1 *   * * *   ubuntu  bash ${HOMETRADE}/main/database/binance/monitor.sh